# config/supervisor/supervisor.conf
[program:scheduler_web]
directory=/home/webmaster/scheduler/
process_name=%(program_name)s
command=/home/webmaster/scheduler/env/bin/gunicorn config.wsgi -b 0.0.0.0:8000
autorestart=true
